# launchd/com.scriptease.jdbc-cli.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scriptease.jdbc-cli</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/java</string>
        <string>-jar</string>
        <string>__JAR_PATH__</string>
        <string>daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/jdbc-cli/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/jdbc-cli/daemon.log</string>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>
</dict>
</plist>
